# docker/supervisor/services/tor.conf
[program:tor]
command = /usr/sbin/tor -f /etc/tor/torrc
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
